# cron.d/gcs-hourly-backup
#!/bin/bash

@hourly gsutil cp /config/World.wld gs://$GCS_BUCKET$GCS_BUCKET_PATH/World.wld
