# cron.d/gcs-hourly-backup
#!/bin/bash

#TODO: implementation and testing of this